# profiles/omarchy/dot_config/systemd/user/kopia-snapshot.timer
[Unit]
Description=Kopia Backup Timer

[Timer]
OnCalendar=hourly
Persistent=true
Unit=kopia-snapshot.service

[Install]
WantedBy=timers.target
