# profiles/omarchy/dot_config/systemd/user/kopia-snapshot.service
[Unit]
Description=Kopia snapshot
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=kopia snapshot create /home/jcontonio
