# profiles/omarchy/dot_config/systemd/user/kopia-snapshot.service
[Unit]
Description=Kopia snapshot → B2
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=kopia snapshot create %2
