# dot_config/autostart/dunst.desktop
[Desktop Entry]
Type=Application
Name=Dunst
Exec=dunst
X-KDE-AutostartPhase=1
X-KDE-StartupNotify=false
NoDisplay=true
